# templates/default/jupyterhub.service.erb
[Unit]
Description=Jupyterhub
After=syslog.target network.target

[Service]
User=<%= node['jupyterhub']['setup']['run_as'] %>
EnvironmentFile=-/etc/sysconfig/jupyterhub
PIDFile=<%= node['jupyterhub']['setup']['pid_file'] %>
<% if node['jupyterhub']['setup']['run_as'] != 'root' -%>
<% if node['jupyterhub']['setup']['enable_ssl'] -%>
ExecStart=/usr/local/bin/jupyterhub -f <%= node['jupyterhub']['setup']['app_dir'] %>/current/jupyterhub_config.py
<% else %>
ExecStart=/usr/local/bin/jupyterhub -f <%= node['jupyterhub']['setup']['app_dir'] %>/current/jupyterhub_config.py --no-ssl
<% end -%>
WorkingDirectory=<%= node['jupyterhub']['setup']['app_dir'] %>/current
<% else -%>
<% if node['jupyterhub']['setup']['enable_ssl'] -%>
ExecStart=/usr/local/bin/jupyterhub -f <%= node['jupyterhub']['setup']['runtime_dir'] %>/jupyterhub_config.py
<% else %>
ExecStart=/usr/local/bin/jupyterhub -f <%= node['jupyterhub']['setup']['runtime_dir'] %>/jupyterhub_config.py --no-ssl
<% end -%>
<% end -%>

[Install]
WantedBy=multi-user.target
